# deployment/sudoku-extraction.service
[Unit]
Description=Sudoku Killer Extraction Service (Python/Flask)
After=network.target

[Service]
# Run as dedicated sudoku user for security
# Create user with: sudo useradd -r -s /usr/sbin/nologin sudoku
User=sudoku
Group=sudoku

# The directory where the Python app is located
WorkingDirectory=/opt/sudoku/extraction_service

# Use gunicorn with single worker (due to model memory requirements)
# Longer timeout for model loading and inference (600s for CPU-based OCR)
ExecStart=/opt/sudoku/extraction_service/venv/bin/gunicorn --bind 127.0.0.1:5001 --workers 1 --timeout 600 app:app

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/sudoku
PrivateTmp=true

Restart=always
RestartSec=5

# Environment
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
