# deployment/sudoku.service
[Unit]
Description=Sudoku Generator Web Server
After=network.target

[Service]
# Load environment variables (e.g. GOOGLE_API_KEY)
EnvironmentFile=/opt/sudoku/config.env

# Replace <user> with your actual username on the server (e.g., ubuntu, root, etc.)
User=root
Group=root

# The directory where the binary and 'dist' folder are located
WorkingDirectory=/opt/sudoku

# The command to start the server
ExecStart=/opt/sudoku/sudoku-server

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
